# com.user.widget-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.widget-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jleos/Library/Mobile Documents/com~apple~ScriptEditor2/Documents/Backup Widgets/backup-widgets.sh</string>
    </array>

    <!-- Run daily at 9:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jleos/Backups/widgets/widget-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jleos/Backups/widgets/widget-backup-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
